# docker/wallet-pawcoin/supervisord.conf
[supervisord]
nodaemon=true

[program:pawcoin]
priority=10
command=/bin/pawcoind -daemon
stdout_logfile=/var/log/supervisor/pawcoin.log
stderr_logfile=/var/log/supervisor/pawcoin.log
autorestart=true